# deployment/origami-backend.service
[Unit]
Description=Origami Importados Backend API
After=network.target postgresql.service

[Service]
Type=notify
User=nginx
Group=nginx
WorkingDirectory=/var/www/origami-backend
ExecStart=/usr/bin/dotnet /var/www/origami-backend/OrigamiBack.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=origami-backend
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# Logs
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
